# tests/regression/debris_weekly_test_schedule.crontab
0 22 * * 0 (. $HOME/.profile ; ./regressionTest ti46)
0 22 * * 1 (. $HOME/.profile ; ./regressionTest tidigits)
0 22 * * 2 (. $HOME/.profile ; ./regressionTest an4_words)
0 22 * * 3 (. $HOME/.profile ; ./regressionTest an4_spelling)
0 22 * * 4 (. $HOME/.profile ; ./regressionTest an4_full)
0 22 * * 5 (. $HOME/.profile ; ./regressionTest an4_plus)
0 22 * * 6 (. $HOME/.profile ; ./regressionTest rm1)
0 20 * * * (. $HOME/.profile ; ./regressionTest standard)
