# tests/regression/filbert_weekly_test_schedule.crontab
MAILTO=cmusphinx-results@lists.sourceforge.net
50 18 * * * (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; svn -q up .; ./regressionTest nightly batch)
35 23 * * 0 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest sunday batch)
35 23 * * 1 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest monday batch)
35 23 * * 2 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest tuesday batch)
35 23 * * 3 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest wednesday batch)
35 23 * * 4 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest thursday batch)
35 23 * * 5 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest friday batch)
35 23 * * 6 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest saturday batch)
50 23 * * 0 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest async0 batch)
50 23 * * 1 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest async1 batch)
50 23 * * 2 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest async2 batch)
50 23 * * 3 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest async3 batch)
50 23 * * 4 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest async4 batch)
50 23 * * 5 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest async5 batch)
50 23 * * 6 (. $HOME/.profile ; cd $SF_ROOT/sphinx4/tests/regression; ./regressionTest async6 batch)
30 17 * * * (. $HOME/.profile ; cd $SF_ROOT/sphinx4/scripts; svn -q up .; ./updateS4Javadocs.sh)
# 05 06 * * * (. $HOME/.profile ; svn -q up $SF_ROOT/sphinx3/src/tests/regression/checkCompileStatus.sh; $SF_ROOT/sphinx3/src/tests/regression/checkCompileStatus.sh)
05 02 * * * (. $HOME/.profile ; cd $SF_ROOT/web; svn -q up .; $SF_ROOT/web/script/nightlyBuild.sh)
05 06 * * * (. $HOME/.profile ; cd $SF_ROOT/web; svn -q up .; $SF_ROOT/web/script/update_sf.sh)
00 03 * * * (. $HOME/.profile ; cd $SF_ROOT/web; svn -q up .; $SF_ROOT/web/script/sfbackup.sh)
