# tests/regression/mangueira_weekly_test_schedule.crontab
0  19 * * * (. $HOME/.profile ; ./regressionTest nightly batch)
45 23 * * 0 (. $HOME/.profile ; ./regressionTest sunday batch)
45 23 * * 1 (. $HOME/.profile ; ./regressionTest monday batch)
45 23 * * 2 (. $HOME/.profile ; ./regressionTest tuesday batch)
45 23 * * 3 (. $HOME/.profile ; ./regressionTest wednesday batch)
45 23 * * 4 (. $HOME/.profile ; ./regressionTest thursday batch)
45 23 * * 5 (. $HOME/.profile ; ./regressionTest friday batch)
45 23 * * 6 (. $HOME/.profile ; ./regressionTest saturday batch)

