# tests/regression/mangueira_weekly_test_schedule.crontab
0 23 * * 0 (. $HOME/.profile ; TEST_NAME=ti46  qsub -v TEST_NAME ./regressionTest )
0 23 * * 1 (. $HOME/.profile ; TEST_NAME=tidigitS  QSub -v TEST_NAME ./regressionTest )
0 23 * * 2 (. $HOME/.profile ; TEST_NAME=an4_wordS  QSub -v TEST_NAME ./regressionTest )
0 23 * * 3 (. $HOME/.profile ; TEST_NAME=an4_spelling  qsub -v TEST_NAME ./regressionTest )
0 23 * * 4 (. $HOME/.profile ; TEST_NAME=an4_full  qsub -v TEST_NAME ./regressionTest )
0 23 * * 5 (. $HOME/.profile ; TEST_NAME=an4_plus  qsub -v TEST_NAME ./regressionTest )
0 23 * * 6 (. $HOME/.profile ; TEST_NAME=rm1  qsub -v TEST_NAME ./regressionTest )
5 11 * * * (. $HOME/.profile ; TEST_NAME=standard  qsub -v TEST_NAME ./regressionTest )
